# config/includes.chroot/etc/profile.d/99-archonsync-umask.sh
# ArchonSync: default to owner-only permissions on newly created files.
umask 077
