# config/includes.chroot/etc/systemd/system/archonsync-windows-mount.service
[Unit]
Description=Mount Windows system partition at /windows
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/archonsync-mount-windows
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
