# config/includes.chroot/etc/systemd/system/archonsync-antivirus.timer
[Unit]
Description=Run the ArchonSync malware scan weekly

[Timer]
OnCalendar=Sun 03:00
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
