# config/includes.chroot/etc/systemd/system/archonsync-blocklist.timer
[Unit]
Description=Weekly ArchonSync blocklist refresh (and shortly after boot)

[Timer]
# Soon after boot, then weekly. Randomised so every install doesn't hammer the
# sources at the same instant.
OnBootSec=3min
OnCalendar=weekly
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
