# config/includes.chroot/etc/xdg/autostart/archonsync-windows-setup.desktop
[Desktop Entry]
Type=Application
Name=ArchonSync Windows Mode Setup
Exec=archonsync-windows-setup
OnlyShowIn=KDE;
X-KDE-autostart-phase=2
