# config/includes.chroot/etc/xdg/autostart/archonsync-setup.desktop
[Desktop Entry]
Type=Application
Name=ArchonSync Setup
Exec=archonsync-setup
Icon=preferences-system
Comment=First-boot setup wizard
X-KDE-autostart-phase=2
OnlyShowIn=KDE;
